# Install/Linux/transcoder.target
[Unit]
After=network.target
Wants=livepeer.service prometheus.service

[Install]
WantedBy=multi-user.target
